# client/com.wireguardvpn.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.wireguardvpn</string>

    <!-- No watch script, no network check: just bring the tunnel up directly. -->
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/wg-quick</string>
        <string>up</string>
        <string>wg0</string>
    </array>

    <!-- wg-quick needs these tools (wg, bash, route, etc.) on PATH under launchd. -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/opt/homebrew/sbin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <!-- Run once at boot/load, unconditionally, on ANY network. -->
    <key>RunAtLoad</key>
    <true/>

    <!-- CRITICAL: wg-quick spawns wireguard-go (the tunnel) then exits. Without this,
         launchd kills wireguard-go when wg-quick exits and the tunnel dies instantly.
         AbandonProcessGroup lets the spawned tunnel process live on. -->
    <key>AbandonProcessGroup</key>
    <true/>

    <!-- NOTE: no WatchPaths (not network-triggered) and no KeepAlive
         (wg-quick exits after the tunnel is up; the interface persists). -->

    <key>StandardOutPath</key>
    <string>/var/log/wireguardvpn.log</string>
    <key>StandardErrorPath</key>
    <string>/var/log/wireguardvpn.log</string>
</dict>
</plist>
